# config/systemd/user/picom.service
[Unit]
Description=Picom

[Service]
ExecStart=/usr/bin/picom --config %h/.config/picom/picom.conf

[Install]
WantedBy=default.target
